# leam/com.greenreach.leam.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.greenreach.leam</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Volumes/CodeVault/Projects/Light-Engine-Foxtrot/leam/leam-agent.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Volumes/CodeVault/Projects/Light-Engine-Foxtrot/leam</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>GREENREACH_API_KEY</key>
        <string>46cd098b9fd96b385546edaf31116876259fd763ebc13d3d5e6dbb36a0e4bc85</string>
        <key>FARM_ID</key>
        <string>FARM-MLTP9LVH-B0B85039</string>
        <key>CENTRAL_WS_URL</key>
        <string>wss://www.greenreachgreens.com</string>
        <key>LEAM_LOG_LEVEL</key>
        <string>info</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/leam-agent.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/leam-agent.log</string>
</dict>
</plist>
