# edge-deployment/esp32-sensor.service
[Unit]
Description=ESP32 Environmental Sensor Reader
After=network.target
Wants=network.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/light-engine
ExecStart=/home/admin/light-engine/venv/bin/python3 /home/admin/light-engine/scripts/esp32-to-env.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Ensure USB serial port is available
# Add user to dialout group: sudo usermod -a -G dialout admin

[Install]
WantedBy=multi-user.target
